# queries/linux_crontab_modification.json
{
  "title": "Linux Crontab Modification",
  "description": "Detects crontab command execution, which may indicate persistence via scheduled tasks.",
  "query": "('Log Source' = 'SOC Linux Syslog Logs' or 'Log Source' = 'Linux Secure Logs' or 'Log Source' = 'Linux Syslog Logs' or 'Log Source' = 'Linux Audit Logs') and (('Process Name' = 'crontab'))",
  "sigma_id": "07b8c9d0-bcde-f012-3456-789012345067",
  "level": "medium",
  "stig_category": "CAT II",
  "tags": [
    "attack.persistence",
    "attack.t1053.003"
  ],
  "mitre_attack": {
    "tactics": [
      "persistence"
    ],
    "techniques": [
      "T1053.003"
    ]
  },
  "logsource": {
    "product": "linux",
    "service": "audit",
    "candidates": [
      "SOC Linux Syslog Logs",
      "Linux Secure Logs",
      "Linux Syslog Logs",
      "Linux Audit Logs"
    ]
  }
}