# queries/windows_registry_run_key_modification.json
{
  "title": "Windows Registry Run Key Modification",
  "description": "Detects modification of Windows Registry Run/RunOnce keys for persistence, allowing malicious programs to execute at system startup or user login.",
  "query": "('Log Source' = 'Windows Sysmon Events' or 'Log Source' = 'Windows Sysmon Operational Logs' or 'Log Source' = 'SOC Windows Sysmon Logs') and ((('Command Line' like '*CurrentVersion\\Run*' or 'Command Line' like '*CurrentVersion\\RunOnce*' or 'Command Line' like '*CurrentVersion\\RunServices*' or 'Command Line' like '*CurrentVersion\\Policies\\Explorer\\Run*')))",
  "sigma_id": "d1a2b3c4-3004-4000-8000-000000000004",
  "level": "high",
  "stig_category": "CAT I",
  "tags": [
    "attack.persistence",
    "attack.t1547.001"
  ],
  "mitre_attack": {
    "tactics": [
      "persistence"
    ],
    "techniques": [
      "T1547.001"
    ]
  },
  "logsource": {
    "product": "windows",
    "service": "process_creation",
    "candidates": [
      "Windows Sysmon Events",
      "Windows Sysmon Operational Logs",
      "SOC Windows Sysmon Logs"
    ]
  }
}