# rules/windows/process_creation/win_scheduled_task_creation.yaml
title: Windows Scheduled Task Creation via Schtasks
id: e1f2a3b4-4018-4000-8000-000000000018
status: stable
description: Detects creation of scheduled tasks via schtasks.exe with suspicious
  patterns. Scheduled tasks are a common persistence and execution mechanism.
logsource:
  product: windows
  service: process_creation
detection:
  selection:
    process_name|endswith: schtasks.exe
    CommandLine|contains:
    - /create
    - -create
  condition: selection
level: medium
tags:
- attack.persistence
- attack.t1053.005
falsepositives:
- Legitimate task scheduling by system administrators
- Software installation creating maintenance tasks
version: 1
